# ansible/roles/middleware/supervisor/files/supervisord.conf
description "supervisord"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
script
   exec /usr/bin/supervisord -n -c /etc/supervisord.conf
end script
